# remount/remount.service
[Unit]
Description=Remount unmounted volumes

[Service]
Type=simple
ExecStart=/root/remount.sh
Restart=no
