# systemd/s3-backup/home-server-backup.timer
[Unit]
Description=Backup home server files to S3 once per day

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
